# web/deploy/systemd/sepsis-dx.service
[Unit]
Description=SepsisDx Web Application (gunicorn)
After=network.target

[Service]
Type=notify
User=sepsis-dx
Group=sepsis-dx
RuntimeDirectory=sepsis-dx
WorkingDirectory=/opt/sepsis-dx/web
Environment="PATH=/opt/sepsis-dx/.venv/bin:/usr/local/bin:/usr/bin"
EnvironmentFile=-/etc/sepsis-dx/environment
ExecStart=/opt/sepsis-dx/.venv/bin/gunicorn --config /opt/sepsis-dx/web/gunicorn.conf.py wsgi:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=10
TimeoutStopSec=30
KillMode=mixed

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/sepsis-dx/models /opt/sepsis-dx/data /opt/sepsis-dx/web/uploads /var/log/sepsis-dx
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
