# web/deploy/systemd/sepsis-dx.service
[Unit]
Description=SepsisDx Web Application (sepsis-dx.com)
After=network.target

[Service]
Type=simple
User=sepsis-dx
Group=sepsis-dx
WorkingDirectory=/opt/sepsis-dx/web

EnvironmentFile=/etc/sepsis-dx/environment

ExecStart=/opt/sepsis-dx/.venv/bin/gunicorn \
    --config /opt/sepsis-dx/web/gunicorn.conf.py \
    wsgi:app

Restart=always
RestartSec=5

ProtectSystem=strict
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/opt/sepsis-dx/models \
               /opt/sepsis-dx/data \
               /opt/sepsis-dx/web/uploads \
               /var/log/sepsis-dx

StandardOutput=journal
StandardError=journal
SyslogIdentifier=sepsis-dx

[Install]
WantedBy=multi-user.target
